# zebra-api.service
[Unit]
Description=Zebra ZPL Print API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/zebra-wifi-bridge
ExecStart=/usr/bin/python3 /home/pi/zebra-wifi-bridge/app.py
Environment=ZPL_API_KEY=changeme
Environment=ZPL_PRINTER_DEVICE=/dev/usb/lp0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
